# firstboot/MicroOS-firstboot.service.in
[Unit]
Description=openSUSE MicroOS Firstboot Service
DefaultDependencies=no
Conflicts=shutdown.target
After=systemd-remount-fs.service
Before=systemd-sysusers.service sysinit.target shutdown.target
ConditionPathIsReadWrite=/etc
ConditionFirstBoot=yes

[Service]
Environment=SYSTEMCTL_OPTIONS=--ignore-dependencies TERM=linux
Type=oneshot
RemainAfterExit=yes
ExecStart=@libexecdir@/MicroOS-firstboot
StandardOutput=tty
StandardInput=tty
StandardError=tty
